# docs/renote-health.service
# /etc/systemd/system/renote-health.service
[Unit]
Description=Renote health check (journal only)
After=network-online.target

[Service]
Type=oneshot
# -L follow redirects, -m 5 timeout, --retry with small delay, ok to log JSON
ExecStart=/usr/bin/curl -LfsS -m 5 --retry 2 --retry-delay 1 --retry-connrefused https://n.cj36.org/api.php?action=health
StandardOutput=journal
StandardError=journal
